# release/systemd/commandcode-router.user.service
[Unit]
Description=CommandCode Router - OpenAI-compatible backend router for CommandCode bridges
After=default.target commandcode-bridge.service
Wants=commandcode-bridge.service

[Service]
Type=simple
WorkingDirectory=%h
EnvironmentFile=%h/.config/commandcode-bridge/env
EnvironmentFile=%h/.config/commandcode-bridge/router.env
ExecStart=%h/.local/bin/commandcode-router
Restart=always
RestartSec=5
TimeoutStopSec=20
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=default.target
