# release/systemd/commandcode-bridge.service
[Unit]
Description=CommandCode Bridge - OpenAI-compatible API for CommandCode DeepSeek
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=commandcode-bridge
Group=commandcode-bridge
WorkingDirectory=/opt/commandcode-bridge
EnvironmentFile=/etc/commandcode-bridge.env
ExecStart=/usr/bin/node /opt/commandcode-bridge/dist/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp
CapabilityBoundingSet=
LockPersonality=true
# MemoryDenyWriteExecute can break Node/V8 JIT on some distributions; keep it disabled unless verified with this exact Node build.

[Install]
WantedBy=multi-user.target
